# .config/systemd/user/batsignal.service
[Unit]
Description=Lightweight battery daemon

[Service]
ExecStart=/usr/bin/batsignal -c 10 -w 30 -n BAT0

[Install]
WantedBy=default.target