# .config/systemd/user/rclone.service
[Unit]
Description=rclone backup

[Service]
Type=oneshot
ExecStart=~/.dotfiles/bin/rclone.sh
Environment="DEST=backblaze:maxholman-backups"
