# bin/90-autorandr.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/mholman/.dotfiles/bin/autorandr.sh"
